# gunicorn.conf
[program:gunicorn]
command=gunicorn --bind localhost:5000 wsgi
directory=/deploy/app